# device-software/meta-edison-middleware/recipes-connectivity/mdns/files/mdns.service
[Unit]
Description=Zero-configuration networking
After=network.target

[Service]
ExecStart=/usr/sbin/mdnsd
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
